# cmd/httpadapter/bin/httpadapter-tunnel.service
[Unit]
Description=HTTP Adapter Tunnel
After=network.target
 
[Service]
Type=simple
WorkingDirectory=/opt/httpadapter
ExecStart=/opt/httpadapter tunnel
KillMode=control-group
Restart=on-failure
LimitNOFILE=100000
 
[Install]
WantedBy=multi-user.target